# rsync-duplicate.service
[Unit]
Description= rsync-duplicate service

[Service]
ExecStart=/usr/local/bin/rsync-duplicate.sh

[Install]
WantedBy=multi-user.target
